# confs/supervisord/apps.conf
[program:gunicorn]
directory = /opt/webapps/signeasy/
command = /bin/bash /startapp.sh
user=root
autostart=true
autorestart=true